# as-ma/systemd/ioc-bo-ma-corrector.service
[Unit]
Description=Soft IOC for BO corrector magnets
After=multi-user.target

[Service]
Group=fac
User=fac
Type=simple
WorkingDirectory=/home/fac/sirius-iocs
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /tmp/bashrc-sirius-ioc-bo-ma-corrector'
EnvironmentFile=-/tmp/bashrc-sirius-ioc-bo-ma-corrector
ExecStart=/usr/local/bin/sirius-ioc-bo-ma-corrector.py
