# as-ma/systemd/ioc-si-ma-multipole-fam.service
[Unit]
Description=Soft IOC for SI dipole magnet families
After=multi-user.target

[Service]
Group=fac
User=fac
Type=simple
WorkingDirectory=/home/fac/sirius-iocs/as-ma
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /tmp/bashrc-sirius-ioc-si-ma-multipole-fam'
EnvironmentFile=-/tmp/bashrc-sirius-ioc-si-ma-multipole-fam
ExecStart=/usr/local/bin/sirius-ioc-si-ma-multipole-fam.py
