# as-ma/systemd/sirius-ioc-bo-ma-quadrupoles.service
[Unit]
Description=Soft IOC for BO quadrupoles
After=multi-user.target

[Service]
Group=sirius
User=sirius
Type=simple
WorkingDirectory=/home/sirius/iocs
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /tmp/bashrc-sirius-ioc-tb-ma-quadrupoles'
EnvironmentFile=-/tmp/bashrc-sirius-ioc-tb-ma-quadrupoles
ExecStart=/usr/local/bin/sirius-ioc-as-ma.py \
BO-Fam:MA-QF \
BO-Fam:MA-QD
