# as-ma/systemd/sirius-ioc-bo-ma-correctors-ia13.service
[Unit]
Description=Soft IOC for BO corrector magnets ia13
After=multi-user.target

[Service]
Group=sirius
User=sirius
Type=simple
WorkingDirectory=/home/sirius/iocs-log
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /home/sirius/iocs-log/env-sirius-ioc-tb-ma-correctors-ia13'
EnvironmentFile=-/home/sirius/iocs-log/env-sirius-ioc-tb-ma-correctors-ia13
ExecStart=/usr/local/bin/sirius-ioc-as-ma.py \
BO-31U:MA-CH \
BO-31U:MA-CV \
BO-33U:MA-CH \
BO-33U:MA-CV
