# as-ma/systemd/sirius-ioc-si-ma-correctors-ia14-m12.service
[Unit]
Description=Power Supply IOC - SI Correctors Subsections 14M1, 14M2
After=multi-user.target

[Service]
Group=sirius
User=sirius
Type=simple
WorkingDirectory=/home/sirius/iocs-log
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /home/sirius/iocs-log/env-sirius-ioc-si-ma-correctors-ia14-m12'
EnvironmentFile=-/home/sirius/iocs-log/env-sirius-ioc-si-ma-correctors-ia14-m12
ExecStart=/usr/local/bin/sirius-ioc-as-ma.py \
SI-14M1:MA-CH \
SI-14M1:MA-CV \
SI-14M2:MA-CH \
SI-14M2:MA-CV
