# as-ma/systemd/sirius-ioc-bo-ma-correctors-ia10.service
[Unit]
Description=Soft IOC for BO corrector magnets ia10
After=multi-user.target

[Service]
Group=sirius
User=sirius
Type=simple
WorkingDirectory=/home/sirius/iocs-log
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/bash --login -c 'env > /home/sirius/iocs-log/env-sirius-ioc-tb-ma-correctors-ia10'
EnvironmentFile=-/home/sirius/iocs-log/env-sirius-ioc-tb-ma-correctors-ia10
ExecStart=/usr/local/bin/sirius-ioc-as-ma.py \
BO-23U:MA-CH \
BO-23U:MA-CV \
BO-25U:MA-CH \
BO-25U:MA-CV
